# templates/default_boot.service
[Unit]
Description=[[DESCRIPTION]]
After=network-online.target

[Service]
Type=oneshot
ExecStart=[[COMMAND]]

[Install]
WantedBy=default.target
