# templates/default_boot.timer
[Unit]
Description=[[DESCRIPTION]]
Requires=[[UNIT_NAME]].service

[Timer]
OnBootSec=10s
RemainAfterElapse=no

[Install]
WantedBy=timers.target

